# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor/
user=genesis

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700
username=genesis
password=believers_commons_genesis

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# Genesis Stack Core Services
[program:genesis-api]
command=node /app/src/genesis-stack.js
directory=/app
user=genesis
environment=NODE_ENV=production,GENESIS_MODE=empire,REVENUE_TARGET=370000000
autostart=true
autorestart=true
stderr_logfile=/var/log/genesis/genesis-api.error.log
stdout_logfile=/var/log/genesis/genesis-api.out.log
redirect_stderr=false
priority=100

[program:riveros-transaction]
command=node /app/src/genesis-stack.js
directory=/app
user=genesis
environment=NODE_ENV=production,GENESIS_MODE=river,REVENUE_SHARE_RATIO=85:15
autostart=true
autorestart=true
stderr_logfile=/var/log/genesis/riveros.error.log
stdout_logfile=/var/log/genesis/riveros.out.log
priority=101

[program:synergizeos-orchestration]
command=node /app/src/genesis-stack.js
directory=/app
user=genesis
environment=NODE_ENV=production,GENESIS_MODE=synergize,COMPLIANCE_AUTOMATION=enabled
autostart=true
autorestart=true
stderr_logfile=/var/log/genesis/synergizeos.error.log
stdout_logfile=/var/log/genesis/synergizeos.out.log
priority=102

# ECG Compliance Framework
[program:ecg-compliance]
command=python3 /app/src/compliance/compliance_server.py
directory=/app
user=genesis
environment=COMPLIANCE_MODE=ecg_charter,ARTICLE_2_1_SECURITY=enabled,ANNUAL_FEE=3000000
autostart=true
autorestart=true
stderr_logfile=/var/log/compliance/ecg-compliance.error.log
stdout_logfile=/var/log/compliance/ecg-compliance.out.log
priority=103

# MAHDI Ruby Sandbox
[program:mahdi-playground]
command=bundle exec rails server -p 3000 -e production
directory=/app/mahdi
user=genesis
environment=RAILS_ENV=production,RUBY_MONETIZATION=premium,OPENTELEMETRY_ENABLED=true
autostart=true
autorestart=true
stderr_logfile=/var/log/mahdi/playground.error.log
stdout_logfile=/var/log/mahdi/playground.out.log
priority=104

# Revenue Optimization Service
[program:revenue-optimizer]
command=python3 /app/src/revenue/optimizer.py
directory=/app
user=genesis
environment=REVENUE_TARGET=370000000,OPTIMIZATION_MODE=aggressive
autostart=true
autorestart=true
stderr_logfile=/var/log/genesis/revenue-optimizer.error.log
stdout_logfile=/var/log/genesis/revenue-optimizer.out.log
priority=105

# Multi-signature Safe Monitor
[program:multi-sig-monitor]
command=python3 /app/src/blockchain/safe_monitor.py
directory=/app
user=genesis
environment=SAFE_THRESHOLD=3,SAFE_OWNERS=5,NETWORK=mainnet
autostart=true
autorestart=true
stderr_logfile=/var/log/genesis/multi-sig.error.log
stdout_logfile=/var/log/genesis/multi-sig.out.log
priority=106

# Compliance Monitor
[program:compliance-monitor]
command=python3 /app/src/compliance/monitor.py
directory=/app
user=genesis
environment=ECG_CHARTER=enabled,DPDPA_2023=enabled,AUTO_REPORTING=true
autostart=true
autorestart=true
stderr_logfile=/var/log/compliance/monitor.error.log
stdout_logfile=/var/log/compliance/monitor.out.log
priority=107

# Cloudflare Security Integration
[program:cloudflare-security]
command=python3 /app/src/security/cloudflare_integration.py
directory=/app
user=genesis
environment=SECURITY_PREMIUM_FEE=4500000,VULNERABILITY_SCANNING=enabled
autostart=true
autorestart=true
stderr_logfile=/var/log/genesis/cloudflare-security.error.log
stdout_logfile=/var/log/genesis/cloudflare-security.out.log
priority=108

# OpenTelemetry Collector
[program:otel-collector]
command=/usr/local/bin/otelcol --config=/app/config/otel-collector.yaml
directory=/app
user=genesis
autostart=true
autorestart=true
stderr_logfile=/var/log/genesis/otel-collector.error.log
stdout_logfile=/var/log/genesis/otel-collector.out.log
priority=109

# NGINX Reverse Proxy (if not using external container)
[program:nginx]
command=nginx -g "daemon off;"
directory=/
user=root
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx/error.log
stdout_logfile=/var/log/nginx/access.log
priority=200

# Group configurations for easier management
[group:genesis-stack]
programs=genesis-api,riveros-transaction,synergizeos-orchestration
priority=100

[group:compliance]
programs=ecg-compliance,compliance-monitor
priority=103

[group:revenue-optimization]
programs=revenue-optimizer
priority=105

[group:blockchain]
programs=multi-sig-monitor
priority=106

[group:security]
programs=cloudflare-security
priority=108

[group:monitoring]
programs=otel-collector
priority=109

[group:infrastructure]
programs=nginx
priority=200